# server/faucet.service
# /etc/systemd/system/faucet.service
[Unit]
Description=Janction Faucet
Wants=network-online.target ipfs.service docker.service
After=network-online.target ipfs.service docker.service

[Service]
User=ec2-user
WorkingDirectory=/home/ec2-user
ExecStart=/usr/local/bin/janction-faucet
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
# Optional hardening / niceties:
NoNewPrivileges=yes
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
